# Linux/installation/pymqtt-listener.service
[Unit]
Description=PyMQTT Topic Listener
After=network.target mosquitto.service
Wants=network.target

[Service]
ExecStart=/usr/bin/python3 /usr/local/bin/pymqtt-listener.py
Restart=always
RestartSec=5
User=root
WorkingDirectory=/usr/local/bin
StandardOutput=file:/var/log/pymqtt-listener.log
StandardError=file:/var/log/pymqtt-listener.log

[Install]
WantedBy=multi-user.target