# Desalinizador/pymqtt-listener.service
[Unit]
Description=PyMQTT Topic Listener
After=network.target mosquitto.service
Wants=network.target

[Service]
ExecStart=/usr/bin/python3 /usr/local/bin/mqtt_listener.py
Restart=always
RestartSec=5
User=root
WorkingDirectory=/usr/local/bin
StandardOutput=append:/var/log/pymqtt-listener.log
StandardError=append:/var/log/pymqtt-listener.log

[Install]
WantedBy=multi-user.target